# worker/configs/worker.supervisor.conf
[supervisord]
nodaemon=true

[program:worker]
command=bash -c "java -jar worker.jar"
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log