# wiki/configs/wiki.supervisor.conf
[supervisord]
nodaemon=true

[program:wiki]
command=bash -c "java -jar wiki.jar"
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log